# hadoop-cdh5/files/datanode.in
#!/bin/sh
#
# PROVIDE: datanode
# REQUIRE: LOGIN
# KEYWORD: nojail

. /etc/rc.subr

name=datanode
rcvar=datanode_enable

export HADOOP_CONF_DIR=%%HADOOP_CONF%%
export HADOOP_PID_DIR=%%HADOOP_RUNDIR%%
export HADOOP_IDENT_STRING=%%HADOOP_USER%%

load_rc_config "${name}"

: ${datanode_enable:=NO}
: ${datanode_user:=%%HADOOP_USER%%}
: ${datanode_group:=%%HADOOP_GROUP%%}

command="/usr/local/sbin/hadoop-daemon.sh"
command_args='start datanode'

stop_cmd=datanode_stop

datanode_stop ()
{
        su -m ${datanode_user} -c "${command} stop datanode"
}

run_rc_command "$1"
